# var/spool/cron/crontabs/root
*/5 * * * * /usr/sbin/logrotate /etc/logrotate.conf
*/27 * * * * /usr/bin/cleanup_var_run.sh
*/15 * * * * /usr/local/cmds/cmdSnap -n1 -rLo
*/6 * * * * /usr/local/cmds/cmdSnap -p
